# deploy/global/tmp/192.168.112.33_16789/supervisor/node-16789.conf
[program:node-16789]
command=/home/lvxiaoyi/./lvxiaoyi_test/node-16789/platon --identity platon --datadir /home/lvxiaoyi/./lvxiaoyi_test/node-16789/data --port 16789 --syncmode 'full' --debug --verbosity 5 --rpc --rpcaddr 0.0.0.0 --rpcport 6789 --rpcapi platon,debug,personal,admin,net,web3 --txpool.nolocals --gcmode archive --nodekey /home/lvxiaoyi/./lvxiaoyi_test/node-16789/data/nodekey --cbft.blskey /home/lvxiaoyi/./lvxiaoyi_test/node-16789/data/blskey
numprocs=1
autostart=false
startsecs=3
startretries=3
autorestart=unexpected
exitcode=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/home/lvxiaoyi/./lvxiaoyi_test/node-16789/log/platon.log
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=20
